# AIM/deploy/systemd/aim-daily-brief-rust.service
[Unit]
Description=AIM daily project brief — pure-Rust path (no Python)
Documentation=file:///home/oem/Desktop/LongevityCommon/AIM/AUDIT_PROJECT_MANAGER_2026-05-06.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.aim_env
WorkingDirectory=%h/Desktop/LongevityCommon/AIM
ExecStart=%h/Desktop/LongevityCommon/AIM/rust-core/target/release/aim-daily-brief --telegram
StandardOutput=journal
StandardError=journal

# Note: alternative to aim-daily-brief.service (which uses Python +
# notify multiplexer). Requires `cargo build --release -p aim-daily-brief`
# beforehand. Both should NOT be enabled simultaneously.
