# deploy/systemd/longevitycommon-autopull.service
[Unit]
Description=LongevityCommon — git pull --ff-only on /home/jaba/web/longevitycommon and trigger restarts on new commits
Documentation=file:///home/oem/Desktop/LongevityCommon/docs/DEPLOY_CONVENTION.md
After=network-online.target

[Service]
Type=oneshot
User=jaba
Group=jaba
WorkingDirectory=/home/jaba/web/longevitycommon
ExecStart=/bin/bash -lc '\
  prev=$(git rev-parse HEAD); \
  git pull --ff-only 2>&1 | head -5; \
  curr=$(git rev-parse HEAD); \
  if [ "$prev" = "$curr" ]; then exit 0; fi; \
  echo "auto-pull: $prev → $curr"; \
  diff_files=$(git diff --name-only "$prev" "$curr"); \
  echo "$diff_files" | head; \
  if echo "$diff_files" | grep -qE "^server/src/|^server/Cargo|^server/build\\.rs|^server/migrations/|^server/\\.sqlx/"; then \
    cd "/home/jaba/web/longevitycommon/server"; \
    echo "[autopull] cargo build --release"; \
    cargo build --release > /tmp/server-rebuild.log 2>&1 && \
      sudo -n systemctl restart longevitycommon-server.service || \
      echo "[autopull] server rebuild FAILED — see /tmp/server-rebuild.log"; \
    cd "/home/jaba/web/longevitycommon"; \
  fi; \
  if echo "$diff_files" | grep -qE "^realtime/"; then \
    cd "/home/jaba/web/longevitycommon/realtime"; \
    . ~/.asdf/asdf.sh; \
    . /etc/aim/lc_social.env; \
    echo "[autopull] mix release"; \
    MIX_ENV=prod mix release --overwrite > /tmp/realtime-rebuild.log 2>&1 && \
      sudo -n systemctl restart longevitycommon-realtime.service || \
      echo "[autopull] realtime rebuild FAILED — see /tmp/realtime-rebuild.log"; \
    cd "/home/jaba/web/longevitycommon"; \
  fi; \
  if echo "$diff_files" | grep -qE "^BioSense/backend/"; then \
    cd "/home/jaba/web/longevitycommon/BioSense/backend"; \
    echo "[autopull] biosense backend cargo build --release"; \
    cargo build --release > /tmp/biosense-rebuild.log 2>&1 && \
      sudo -n systemctl restart biosense-backend.service || \
      echo "[autopull] biosense rebuild FAILED — see /tmp/biosense-rebuild.log"; \
    cd "/home/jaba/web/longevitycommon"; \
  fi; \
  if echo "$diff_files" | grep -qE "^web/"; then \
    cd "/home/jaba/web/longevitycommon/web"; \
    echo "[autopull] npm install + build"; \
    (npm install --no-audit --no-fund > /tmp/web-rebuild.log 2>&1 && \
     npm run build >> /tmp/web-rebuild.log 2>&1 && \
     sudo -n rsync -a --delete dist/ /var/www/longevitycommon-web/) || \
      echo "[autopull] web rebuild FAILED — see /tmp/web-rebuild.log"; \
    cd "/home/jaba/web/longevitycommon"; \
  fi; \
  if echo "$diff_files" | grep -qE "^deploy/nginx/"; then \
    echo "[autopull] nginx site refresh"; \
    sudo -n install -m 0644 deploy/nginx/app.longevity.ge.conf /etc/nginx/sites-available/app.longevity.ge.conf && \
      sudo -n nginx -t && sudo -n systemctl reload nginx || \
      echo "[autopull] nginx reload FAILED"; \
  fi; \
  exit 0'
StandardOutput=journal
StandardError=journal
