# deploy/systemd/ze-web.service
[Unit]
Description=Ze Web — Phoenix LiveView simulator (native)
After=network.target
Wants=network.target

[Service]
Type=exec
User=jaba
Group=jaba
WorkingDirectory=/opt/ze-web

# Pin to host loopback port that nginx already proxies (4400).
Environment="PORT=4400"
Environment="PHX_HOST=ze.longevity.ge"
Environment="PHX_SERVER=true"
Environment="LANG=en_US.UTF-8"
Environment="LC_ALL=en_US.UTF-8"
Environment="ZE_BACKEND_URL=http://127.0.0.1:4401"
Environment="SECRET_KEY_BASE=CM7l1JMUSlCZwLGbDn_j_OKtyEfRt27eZLgOXpxx-xxZypzXQqJzfIv8sqLYxEkA"
Environment="LONGEVITYCOMMON_VERSION=v5.6"
Environment="RELEASE_COOKIE=ze-web-prod-cookie"

ExecStart=/opt/ze-web/bin/ze_web start
ExecStop=/opt/ze-web/bin/ze_web stop
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/ze-web/tmp /tmp

[Install]
WantedBy=multi-user.target
